# lights/crontab
00 00 * * * ~/workspace/automation/lights/schedlights.sh
