# lights/crontab
00 00 * * * ~/workspace/automation/lights/schedlights.sh > /dev/null 2>&1
